# esp8266_server.conf
[program:esp8266]
directory=/root/esp8266_iot_node/
command=python /root/esp8266_iot_node/server.py
autorestart=true
user=root
redirect_stderr=true
stdout_logfile=/root/supervisor_log/log.txt
stderr_logfile=/root/supervisor_log/err.txt
environment=PATH="/opt/xtensa-lx106-elf/bin:%(ENV_PATH)s"